# crontab
# Apagar el bot a las 12:30 AM Cuba (UTC-4)
30 4 * * * pkill -f "python bot.py"

# Encender el bot a las 7:49 AM Cuba (UTC-4)
49 11 * * * cd /app && python bot.py &
